# infrastructure/systemd/morning-signal.service.d/10-after-news.conf
# Order the podcast AFTER the daily-news collector so the generated episode can
# inject the freshest pre-fetched news digest. Wants= (not Requires=): the news
# digest is SUPPLEMENTARY — morning-signal does its own live web search and can
# fall back to the last digest — so a daily-news FAILURE must degrade the digest,
# NOT cancel the whole episode (which Requires= would do). See morning-signal#78.
#
# NB: After= still makes the podcast WAIT for daily-news to finish; for fast,
# quota-safe manual recovery that skips a fresh news sweep entirely, use
# morning-signal-recover.sh (generate-only against the existing digest).
[Unit]
Wants=daily-news.service
After=daily-news.service
