# infrastructure/systemd/morning-signal.service.d/10-memory.conf
[Service]
MemoryAccounting=yes
MemoryHigh=600M
MemoryMax=900M
