# infrastructure/systemd/morning-signal-bakeoff.service
[Unit]
Description=Morning Signal OSS Phase B shadow bakeoff (config#1659) — weekly candidate-vs-prod comparison, never published
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/morning-signal
Environment="MORNING_SIGNAL_RUNNER_ROLE_ARN=arn:aws:iam::711398986525:role/morning-signal-runner-role"
Environment="MORNING_SIGNAL_USE_SSM=1"
Environment="MORNING_SIGNAL_SSM_REGION=us-east-1"
Environment="MORNING_SIGNAL_BAKEOFF_LOG_DIR=/home/ec2-user/morning-signal/bakeoff_logs"
Environment="PATH=/home/ec2-user/morning-signal/.venv/bin:/usr/local/bin:/usr/bin:/bin"
# Refresh the codebase before every run (coupled to the run, not a clock) so
# any merged PR to oss_bakeoff.py/claude.py is live on the next comparison.
# Best-effort (leading '-'): a transient git/pip blip must never SKIP a
# scheduled comparison -- it runs last-good code (mirrors morning-signal.service).
ExecStartPre=-/usr/bin/git fetch origin --quiet
ExecStartPre=-/usr/bin/git reset --hard origin/main
ExecStartPre=-/home/ec2-user/morning-signal/.venv/bin/pip install -e . --quiet
# No --date/--edition: defaults to today's date + the "am" edition, comparing
# against whatever the real production episode used today.
ExecStart=/home/ec2-user/morning-signal/.venv/bin/python scripts/oss_bakeoff.py
StandardOutput=journal
StandardError=journal
TimeoutStartSec=900
PrivateTmp=true

[Install]
WantedBy=multi-user.target
