# infrastructure/systemd/substrate-health-daily.service
[Unit]
Description=Daily-cadence transparency substrate health check (re-homed from ne-postclose-trading-pipeline's DailySubstrateHealthCheck chain, alpha-engine-config-I2722)
Wants=network-online.target
After=network-online.target
# config#2954: a failed/never-finalized run was previously invisible — no
# alert fed off this unit's own state. OnFailure= fires the generic
# alert-on-failure@ handler (alert_on_failure.sh) via krepis.alerts,
# independent of whatever this script's own run does or doesn't manage to
# do before dying.
OnFailure=alert-on-failure@%n.service

[Service]
Type=oneshot
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/alpha-engine-dashboard
ExecStart=/home/ec2-user/alpha-engine-dashboard/infrastructure/substrate_health_check_daily.sh
StandardOutput=journal
StandardError=journal
# Mirrors the SF's executionTimeout=180s script budget; extra headroom since
# there's no separate SSM delivery-timeout layer bounding this path.
TimeoutStartSec=300
RemainAfterExit=no
# config#2954: /var/log/ itself isn't writable by User=ec2-user. systemd
# creates+owns /var/log/substrate-health-daily/ for this unit's lifetime
# (mode 0755, User:Group owned) before ExecStart runs — no bespoke chown/
# tmpfiles.d entry needed. Requires systemd >= 235 (AL2023 ships 252).
LogsDirectory=substrate-health-daily

[Install]
WantedBy=multi-user.target
